# docs/operations/systemd/radiofy-crawl@.service
[Unit]
Description=Radiofy crawl (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/radiofy
ExecStart=/usr/local/bin/bun run crawl --station=%i

[Install]
WantedBy=default.target
